# linux-mt/rootfs/etc_ro/inittab
::sysinit:/etc_ro/rcS
ttyMT2::respawn:/bin/sh
ttyS2::respawn:/bin/sh
ttyS0::respawn:/bin/sh
